# r4s_files/etc/uci-defaults/99-custom
cat /etc/config/socat-opkg >> /etc/config/socat
uci set luci.main.mediaurlbase=/luci-static/darkmatter
uci commit luci
